# systemd/hotbar.service
[Unit]
Description=Hotbar file history panel
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/hotbar
Restart=on-failure
RestartSec=3
Environment=RUST_LOG=info

[Install]
WantedBy=graphical-session.target
